# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/ash -c "httpd -D FOREGROUND"

[program:php-fpm8.1]
command=/bin/ash -c "php-fpm -F"
